# serverHTTP/var/spool/cron/crontabs/root~
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.la83VE/crontab installed on Fri May 22 08:06:22 2015)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# m h  dom mon dow   command
* * * * * sh /root/backup

